# hardware/dh3/fs/etc/systemd/system/link-speed.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/link-speed.sh

[Install]
WantedBy=default.target